# SORGENTI/overlay/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStartPre=-/usr/local/bin/quelo-setup-display
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
Type=idle
